# systemd/system/snap-google\x2dcloud\x2dsdk-172.mount
[Unit]
Description=Mount unit for google-cloud-sdk, revision 172
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/google-cloud-sdk_172.snap
Where=/snap/google-cloud-sdk/172
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
